# cli/crates/ec-deploy/tests/fixtures/dallas/golden/dallas-console/supervisord.conf
; Generated by ec-deploy (Deployment Studio HOST renderer). Do not edit by hand.

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[unix_http_server]
file=/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[program:emqx]
command=/usr/bin/emqx foreground
user=emqx
environment=HOME="/var/lib/emqx",EMQX_ALLOW_ANONYMOUS="true"
priority=10
autorestart=true
startsecs=10
startretries=5
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:config-component]
command=/usr/local/bin/wait-for-tcp localhost:1883 -- config-component --platform HOST --transport MQTT /config/config-component-messaging.json -c FILE /config/config-component-config.json -t dallas-console
priority=15
autorestart=true
startsecs=5
stopwaitsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:edge-console]
command=/usr/local/bin/wait-for-tcp localhost:1883 -- edge-console-gateway --platform HOST --transport MQTT /config/console-messaging.json -c CONFIG_COMPONENT -t dallas-console
directory=/app/edge-console
environment=MALLOC_ARENA_MAX="2",EDGECONSOLE_WORKER_THREADS="4"
priority=20
autorestart=true
startsecs=5
stopwaitsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
