# templates/java/supervisor/component.conf
; supervisord program block for <<COMPONENTNAME>> on the HOST platform.
;
; NOTE: in a rendered deployment this file is OWNED BY THE DEPLOYMENT RENDERER — the EdgeCommons
; Deployment Studio HOST renderer regenerates it from the deployment model (install paths, config
; paths, thing name, user, log dirs). Hand-edits are overwritten on the next render. Treat the copy
; in this repo as the template default / a standalone reference you can run yourself:
;
;   sudo cp supervisor/component.conf /etc/supervisor/conf.d/<<COMPONENTNAME>>.conf
;   sudo supervisorctl reread && sudo supervisorctl update
;   sudo supervisorctl status <<COMPONENTNAME>>
;
; Expected layout (what the renderer lays down):
;   /opt/edgecommons/<<COMPONENTNAME>>/<<JARNAME>>-1.0.0.jar   the shaded component jar (mvn package)
;   /etc/edgecommons/<<COMPONENTNAME>>/config.json             the component config    (-c FILE)
;   /etc/edgecommons/<<COMPONENTNAME>>/messaging.json          the MQTT broker config  (--transport MQTT)
;   /var/log/edgecommons/                                      stdout/stderr logs

[program:<<COMPONENTNAME>>]
; --enable-native-access=ALL-UNNAMED is required by the telemetry-streaming FFM binding.
command=java --enable-native-access=ALL-UNNAMED -jar /opt/edgecommons/<<COMPONENTNAME>>/<<JARNAME>>-1.0.0.jar --platform HOST --transport MQTT /etc/edgecommons/<<COMPONENTNAME>>/messaging.json -c FILE /etc/edgecommons/<<COMPONENTNAME>>/config.json -t <<BINNAME>>-dev
directory=/opt/edgecommons/<<COMPONENTNAME>>
user=edgecommons

autostart=true
autorestart=true
startsecs=10
startretries=3
priority=200

; The library installs a SIGTERM handler that runs its graceful, idempotent shutdown (flip /readyz
; to 503, unsubscribe, close messaging/metrics/streams). Give the JVM room to finish it.
stopsignal=TERM
stopwaitsecs=30
stopasgroup=true
killasgroup=true

stdout_logfile=/var/log/edgecommons/<<COMPONENTNAME>>.out.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/edgecommons/<<COMPONENTNAME>>.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
redirect_stderr=false
