# cli/crates/ec-deploy/tests/fixtures/dallas/golden/gw-fill-01/supervisord.conf
; Generated by ec-deploy (Deployment Studio HOST renderer). Do not edit by hand.

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[unix_http_server]
file=/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[program:emqx]
command=/usr/bin/emqx foreground
user=emqx
environment=HOME="/var/lib/emqx",EMQX_ALLOW_ANONYMOUS="true"
priority=10
autorestart=true
startsecs=10
startretries=5
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:field-sim]
command=/opt/simenv/bin/python /opt/sims/dallas_filling_sim.py
priority=20
autorestart=true
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:config-component]
command=/usr/local/bin/wait-for-tcp localhost:1883 -- config-component --platform HOST --transport MQTT /config/config-component-messaging.json -c FILE /config/config-component-config.json -t gw-fill-01
priority=25
autorestart=true
startsecs=5
stopwaitsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:uns-bridge]
command=/usr/local/bin/wait-for-tcp localhost:1883 -- /bin/bash -lc 'sleep 3; exec uns-bridge --platform HOST --transport MQTT /config/bridge-messaging.json -c CONFIG_COMPONENT -t gw-fill-01'
priority=28
autorestart=true
startsecs=5
stopwaitsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:opcua-adapter]
command=/usr/local/bin/wait-for-tcp localhost:1883 localhost:4840 -- /bin/bash -lc 'sleep 6; exec java --enable-native-access=ALL-UNNAMED -jar /app/opcua/app.jar --platform HOST --transport MQTT /config/opcua-messaging.json -c CONFIG_COMPONENT -t gw-fill-01'
priority=30
autorestart=true
startsecs=8
stopwaitsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:modbus-adapter]
command=/usr/local/bin/wait-for-tcp localhost:1883 localhost:5020 -- /bin/bash -lc 'sleep 6; exec /opt/pyenv/bin/python /app/modbus/main.py --platform HOST --transport MQTT /config/modbus-messaging.json -c CONFIG_COMPONENT -t gw-fill-01'
directory=/app/modbus
priority=30
autorestart=true
startsecs=5
stopwaitsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:telemetry-processor]
command=/usr/local/bin/wait-for-tcp localhost:1883 -- /bin/bash -lc 'sleep 6; exec telemetry-processor --platform HOST --transport MQTT /config/telemetry-messaging.json -c CONFIG_COMPONENT -t gw-fill-01'
priority=40
autorestart=true
startsecs=5
stopwaitsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:file-replicator]
command=/usr/local/bin/wait-for-tcp localhost:1883 -- /bin/bash -lc 'sleep 6; exec file-replicator --platform HOST --transport MQTT /config/filerep-messaging.json -c CONFIG_COMPONENT -t gw-fill-01'
priority=50
autorestart=true
startsecs=5
stopwaitsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
